# files/etc/init.d/blue-merle-rf-guard
#!/bin/sh /etc/rc.common

START=12
STOP=90
USE_PROCD=1

display_msg() {
	echo "{\"msg\": \"$1\"}" > /dev/ttyS0 2>/dev/null || true
}

force_rf_off_early() {
	local waited=0

	display_msg "RF-guard:       booting..."

	while [ "$waited" -lt 30 ]; do
		if gl_modem AT "AT" 2>/dev/null | grep -q "OK"; then
			if gl_modem AT "AT+CFUN=4" 2>/dev/null | grep -q "OK"; then
				logger -p notice -t blue-merle-rf-guard "Early boot RF-off applied."
				display_msg "RF-guard:       RF disabled."
				return 0
			fi
		fi

		waited=$((waited + 1))
		sleep 1
	done

	logger -p notice -t blue-merle-rf-guard "Early boot RF-off could not be applied within 30s."
	display_msg "RF-guard:       RF-off failed!"
	return 1
}

start_service() {
	[ -f "/proc/gl-hw-info/build-in-modem" ] || return 0
	force_rf_off_early

	procd_open_instance
	procd_set_param command /usr/bin/blue-merle-rf-guard
	procd_set_param respawn 300 5 10
	procd_close_instance
}
